# scripts/macos/com.laskobot.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.laskobot.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__INSTALL_DIR__/dist/index-unified.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__NODE_BIN_DIR__:/usr/local/bin:/usr/bin:/bin</string>
        <key>BROWSER_MCP_WS_PORT</key>
        <string>8765</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/laskobot.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/laskobot.err</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
